# supervisord.conf
[unix_http_server]
file=/var/tmp/data_service.supervior.sock
chmod=0760
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
nodaemon=true
logfile=/var/tmp/data_service_supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
user=root


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl = unix:///var/tmp/data_service.supervior.sock

; use nodaemon to start proc
[program:app]
command=python3 /app/app.py
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup=true
killasgroup=true